# build/cozy/root/etc/systemd/system/k3s-cleanup.service
[Unit]
Description=Cleanup pods terminated by node shutdown
Wants=k3s.service

[Service]
Type=oneshot
Environment=KUBECONFIG=/etc/rancher/k3s/k3s.yaml
ExecStart=kubectl delete pods --field-selector status.phase=Failed -A --ignore-not-found=true
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target