# build/cozy/root/etc/containers/systemd/elasticsearch.container
[Unit]
Description=Elasticsearch for TAO Community
Before=k3s.service

[Container]
Image=docker.io/elasticsearch:8.17.4
#Image=docker.elastic.co/elasticsearch/elasticsearch:8.17.4
IP=10.247.40.7
Network=deps.network
NetworkAlias=es
Environment=node.name=es
Environment=discovery.type=single-node
Environment=ES_JAVA_OPTS="-Xms512m -Xmx512m"
Environment=DISABLE_SECURITY_PLUGIN=true
Environment=xpack.security.enabled=false
Environment=xpack.security.transport.ssl.enabled=false
Environment=xpack.security.http.ssl.enabled=false
Volume=/var/lib/tao-ce/dependencies/es:/usr/share/elasticsearch/data:z,U

[Service]
Restart=always
ExecStartPre=mkdir -p /var/lib/tao-ce/dependencies/es

[Install]
WantedBy=multi-user.target