# build/cozy/root/etc/containers/systemd/firestore-emulator.container
[Unit]
Description=Firestore emulator for TAO Community

[Container]
Image=quay.io/tao-ce/services/firestore-emulator:latest
IP=10.247.40.9
Network=deps.network
NetworkAlias=firestore
Environment=PROJECT_ID=demo-tao
Volume=/var/lib/tao-ce/dependencies/firestore:/data:z,U

[Service]
Restart=always
ExecStartPre=mkdir -p /var/lib/tao-ce/dependencies/firestore
RestartSec=5s

[Install]
WantedBy=multi-user.target