# build/cozy/root/etc/containers/systemd/redis.container
[Unit]
Description=Redis for TAO Community
Before=k3s.service

[Container]
Image=docker.io/valkey/valkey:8-alpine
IP=10.247.40.4
NetworkAlias=redis
Network=deps.network
Volume=/var/lib/tao-ce/dependencies/redis:/data:z,U

[Service]
Restart=always
ExecStartPre=mkdir -p /var/lib/tao-ce/dependencies/redis

[Install]
WantedBy=multi-user.target