# libexec/systemd/tao-ce.target
[Unit]
Description=TAO Community Edition - Monolith

[Install]
WantedBy=multi-user.target